# installers/linux/systemd/novakey.service
[Unit]
Description=NovaKey Secure Typing Service (user session)
After=default.target

[Service]
Type=simple

# --- Hardening (what you had before) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict

# You WRITE into the home directory, so don't block it entirely.
# This makes ~ read-only except for the paths listed below.
ProtectHome=read-only
ReadWritePaths=%h/.local/share/novakey %h/.config/novakey

LockPersonality=true
MemoryDenyWriteExecute=true
UMask=0077

Restart=on-failure
RestartSec=1

# --- Installer-safe startup (no WorkingDirectory directive) ---
# 1) create dirs
# 2) seed config on first run
# 3) cd to data dir (so any relative fallback like server_config.json works)
# 4) exec novakey
ExecStart=/bin/sh -lc '\
  set -e; \
  /bin/mkdir -p "%h/.local/share/novakey/logs" "%h/.config/novakey"; \
  if [ ! -f "%h/.local/share/novakey/server_config.yaml" ]; then \
    /bin/cp /usr/share/novakey/server_config.yaml "%h/.local/share/novakey/server_config.yaml"; \
  fi; \
  cd "%h/.local/share/novakey"; \
  exec /usr/bin/novakey --config "%h/.local/share/novakey/server_config.yaml" \
'

[Install]
WantedBy=default.target

